# systemd/.config/systemd/user/polybar-dp-2.service
[Unit]
Description=polybar on DP-2
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/polybar DP-2

[Install]
WantedBy=i3-session.target
